# install_files/outlook-email-checker.timer
[Unit]
Description=Checks for new emails each minute

[Timer]
OnCalendar=*:0/1
Persistent=true
AccuracySec=30sec

[Install]
WantedBy=timers.target
